# docker/supervisor/supervisor.d/dev/yarn-watch.ini
[program:yarn-watch]
directory=/app
command=yarn encore dev --watch
redirect_stderr=true
stdout_logfile=/var/log/supervisor/yarn-watch.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
