# docker/supervisor/supervisor.d/app-ws-server.ini
[program:app-ws-server]
directory=/srv/app
command=php bin/console app:ws-server

