# docker/supervisor/supervisor.d/app-ws-server.ini
[program:app-ws-server]
directory=/app
command=php bin/console app:ws-server -v
redirect_stderr=true
stdout_logfile=/var/log/supervisor/app-ws-server.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
