# docker/supervisor/supervisor.d/app-apm-listener.ini
[program:apm-listener]
directory=/app
command=php -d memory_limit=-1 bin/console app:apm-listener -v
redirect_stderr=true
stdout_logfile=/var/log/supervisor/apm-listener.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
