# docker/supervisor/supervisor.d/messenger-worker.ini
[program:messenger-consume]
directory=/srv/app
command=php -d memory_limit=-1 bin/console messenger:consume analytics async -vv
user=www-data
numprocs=1
startsecs=0
autostart=true
autorestart=true
startretries=10
process_name=analytics_p%(process_num)02d
redirect_stderr=true
stdout_logfile=/var/log/supervisor/messenger-worker.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
