# docker/supervisor/supervisor.d/monitor-docker-api-events.ini
[program:monitor-docker-api-events]
directory=/app
command=php -d memory_limit=-1 bin/console app:monitor-docker-api-events
redirect_stderr=true
stdout_logfile=/var/log/supervisor/monitor-docker-api-events.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
