# docker/supervisor/supervisor.d/msw.ini
[program:msw]
directory=/app/simulations/debian.12-x64
command=/app/simulations/debian.12-x64/MSW 'Port=%(ENV_WATCHDOG_PORT)s'
redirect_stderr=true
stdout_logfile=/var/log/supervisor/msw.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
