# docker/supervisor/supervisor.d/msw.ini
[program:msw]
directory=/srv/app/simulations/v1
command=/srv/app/simulations/v1/MSW
