# docker/supervisor/supervisor.d/msw.ini
[program:msw]
directory=/srv/app/simulations/v1
command=/srv/app/simulations/v1/MSW 'Port=%(ENV_WATCHDOG_PORT)s'
redirect_stderr=true
stdout_logfile=/var/log/supervisor/msw.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
